# lib/systemd/user/observer/observe-youtube-music-conf.service
[Unit]
Description=observe-youtube-music-conf
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/workspace/cron-conf/raspberrypi/bin
ExecStart=/home/ubuntu/workspace/cron-conf/raspberrypi/bin/observe-dnsmasq-toggle-conf --event create --file youtube_music.conf --target /etc/dnsmasq.toggle --mode unlink
TimeoutSec=15
Restart=always

[Install]
WantedBy=default.target
