# lib/systemd/user/observe-cron-conf.service
[Unit]
Description=observe-cron-conf
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/workspace/cron-conf/raspberrypi/bin
ExecStart=/home/ubuntu/workspace/cron-conf/raspberrypi/bin/observe-dnsmasq-toggle-conf --event modify --target /home/ubuntu/workspace/cron-conf
TimeoutSec=15
Restart=always

[Install]
WantedBy=default.target
